# cron/crontabs/test
* * * * * /scripts/test
* * * * * echo "haha"
# Mandatory blank line, don't use . in the script name
